# deploy/systemd/cloud-cli-proxy-host-agent.service
[Unit]
Description=Cloud CLI Proxy host-agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=cloudproxy
Group=cloudproxy
WorkingDirectory=/opt/cloud-cli-proxy
Environment=HOST_AGENT_SOCKET=/run/cloud-cli-proxy/host-agent.sock
ExecStartPre=/usr/bin/env bash /opt/cloud-cli-proxy/deploy/scripts/host-preflight.sh
ExecStart=/opt/cloud-cli-proxy/bin/host-agent
RuntimeDirectory=cloud-cli-proxy
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN
NoNewPrivileges=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
