# deploy/systemd/cloud-cli-proxy-control-plane.service
[Unit]
Description=Cloud CLI Proxy control-plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloudproxy
Group=cloudproxy
WorkingDirectory=/opt/cloud-cli-proxy
Environment=CONTROL_PLANE_ADDR=:8080
ExecStart=/opt/cloud-cli-proxy/bin/control-plane
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
